# app/configs/utils/supervisord.conf
[supervisord]
nodaemon=true
user=root

[supervisorctl]

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
user=root

[program:cron]
command=/usr/sbin/cron -f
user=root